# systemd/.config/systemd/user/guardian.service
[Unit]
Description=update wallpapers

[Service]
ExecStart=/usr/local/bin/guardian_pod.sh
