# systemd/.config/systemd/user/wall_change.service
[Unit]
Description=change the wallpaper

[Service]
ExecStart=/usr/local/bin/wall_change.sh
